# dot-config/anacron/crontab
@hourly anacron -t "$HOME/.config/anacron/anacrontab" -S "$HOME/.local/share/anacron"
