# .config/anacron/cron.daily/8550-5331_costa-verde
#!/bin/sh
rclone copy $HOME/Documents/8550-5331_CostaVerde/ DriveEncrypted:8550-5331_CostaVerde
