# .config/anacron/cron.daily/origami
#!/bin/sh
rclone sync $HOME/Pictures/Origami Drive:Backups/Origami
