# .config/anacron/cron.daily/pii
#!/bin/sh
rclone copy $HOME/Documents/PII DriveEncrypted:PII
