# .config/anacron/cron.daily/plocate
#!/bin/sh
updatedb -l 0 -o $XDG_DATA_HOME/plocate.db -U $HOME
