# dot-config/anacron/cron.daily/baba-is-you
#!/bin/sh
rclone sync $HOME/.local/share/Baba_Is_You/ Drive:Backups/Baba_Is_You
