# dot-config/anacron/cron.daily/books
#!/bin/sh
rclone copy $HOME/Documents/Books Drive:Backups/Books
