# dot-config/anacron/cron.daily/books
#!/bin/sh
rclone sync $HOME/Documents/Books Drive:Backups/Books
