# dot-config/anacron/cron.daily/create-above-and-beyond
#!/bin/sh
rclone copy $HOME'/.local/share/PrismLauncher/instances/Create- Above and Beyond/minecraft/saves' 'Drive:Backups/Create- Above and Beyond'
