# dot-config/anacron/cron.daily/cross-code
#!/bin/sh
rclone copyto $HOME/.config/CrossCode/Default/cc.save Drive:Backups/CrossCode/cc.save
