# dot-config/anacron/cron.daily/origami
#!/bin/sh
rclone copy $HOME/Pictures/Origami Drive:Backups/Origami
