# dot-config/anacron/cron.daily/unity-3d
#!/bin/sh
rclone copy $HOME/.config/unity3d/ Drive:Backups/unity3d
